# images/kanidm-bunny/rootfs/etc/s6-overlay/s6-rc.d/kanidm-ops-api/run
#!/bin/sh
exec /command/with-contenv /usr/local/bin/run-kanidm-ops-api
