# images/kanidm-bunny/rootfs/etc/s6-overlay/s6-rc.d/kanidm-server/finish
#!/bin/sh
exec /usr/local/bin/finish-kanidm-server "$@"
